# images/files/husk-poweroff.service
[Unit]
Description=Power off the Husk slot after the runner exits (recycle trigger)

[Service]
Type=oneshot
ExecStart=/sbin/poweroff
